# scripts/ogn-decode.socket
[Socket]
ListenFIFO=/var/run/ogn-decode.fifo
SocketUser=ogn

[Install]
WantedBy=sockets.target
